# scripts/systemd/cctv-encoder.service
[Unit]
Description=NVR Cam — AV1 Idle Encoder
After=network.target

[Service]
Type=simple
User=nvr
WorkingDirectory=/opt/nvr_cam
EnvironmentFile=/opt/nvr_cam/.env
ExecStart=/opt/nvr_cam/venv/bin/python -m backend.services.encoder.scheduler
Restart=always
RestartSec=60
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
