# scripts/systemd/cctv-motion.service
[Unit]
Description=NVR Cam — Motion Detection Engine
After=network.target nvr-recorder.service

[Service]
Type=simple
User=nvr
WorkingDirectory=/opt/nvr_cam
EnvironmentFile=/opt/nvr_cam/.env
ExecStart=/opt/nvr_cam/venv/bin/python -m backend.services.motion.detector
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
